# config/systemd/user/system-maintenance.timer
[Unit]
Description=Weekly system maintenance report (updates + orphans)

[Timer]
OnCalendar=Sun 10:00
Persistent=true

[Install]
WantedBy=timers.target
